# payment-deb/src/control/prerm
#!/bin/sh
set -eu
systemctl stop pressreleases.timer
systemctl disable pressreleases.timer
