# payment-deb/src/control/prerm
#!/bin/sh
systemctl stop payment  || true
systemctl disable payment  || true
