# payment-deb/src/systemd/pressreleases.service
[Unit]
Description=Press release importer

[Service]
Type=oneshot
User=pressreleases
Group=pressreleases
ExecStart=/opt/pressreleases/run
